# etc/cron.d/eyepea_monitoring_agent
MAILTO=""
*/5 * * * * root /usr/bin/env eyepea_monitoring_agent &> /dev/null