# etc/cron.d/tanto
MAILTO=""
* * * * * root /usr/bin/env tanto > /dev/null 2>&1
